# s6-overlay/s6-rc.d/fail2ban/run
#!/bin/bash
exec fail2ban-server -xf start 2> >(stderr-to-json)
