# s6-overlay/s6-rc.d/sshd/run
#!/bin/bash
exec /usr/sbin/sshd -D -e 2> >(stderr-to-json)
