# s6-overlay/s6-rc.d/syslog/run
#!/bin/bash
exec rsyslogd -n 2> >(stderr-to-json)
